# cron.d/hsec_secalert
#
# {{ ansible_managed }}
#
# Vuln alerts
#
# NVD
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-nvd.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-nvd.pl  > /dev/null 2>&1

# REDHAT
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-redhat.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-redhat.pl  > /dev/null 2>&1

# CENTOS ERRATA
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-centos-errata.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-centos-errata.pl  > /dev/null 2>&1


# DEBIAN
#55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-debian.sh  > /dev/null 2>&1
#58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-debian.pl  > /dev/null 2>&1

# UBUNTU
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-ubuntu.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-ubuntu.pl  > /dev/null 2>&1

# VMWARE
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-vmware.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-vmware.pl  > /dev/null 2>&1

# MICROSOFT
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-microsoft.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-microsoft.pl  > /dev/null 2>&1

# EXPLOITDB
55 1,5,8,14,19,23 * * * root /bin/bash /localservices/sbin/secalert-exploitdb.sh  > /dev/null 2>&1
58 23 * * * root /usr/bin/perl /localservices/sbin/secalert-exploitdb.pl  > /dev/null 2>&1
